# templates/systemd/filebeat.composite.service.erb
[Unit]
Description=Filebeat composite service to start multiple outputs.
Documentation=https://www.elastic.co/products/beats/filebeat
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
